# cloudflare.ddns.service
[Unit]
Description=Cloudflare DDNS Updater
After=network.target

[Service]
ExecStart=/path/to/script
Restart=always
User=

[Install]
WantedBy=multi-user.target